# templates/default/axonops-search.service.erb
[Unit]
Description=AxonOps Search (Elasticsearch)
Documentation=https://www.elastic.co
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
RuntimeDirectory=axonops-search
PrivateTmp=true
Environment=ES_JAVA_HOME=<%= @java_home %>
Environment=ES_HOME=<%= @elastic_install_dir %>
Environment=ES_PATH_CONF=<%= @elastic_install_dir %>/config
Environment=PID_DIR=/var/run/axonops-search

WorkingDirectory=<%= @elastic_install_dir %>

User=<%= @elastic_user %>
Group=<%= @elastic_group %>

ExecStart=<%= @elastic_install_dir %>/bin/elasticsearch

# StandardOutput is configured to redirect to journalctl since
# some systemd versions do not show logs when using 'append'
StandardOutput=journal
StandardError=inherit

# Specifies the maximum file descriptor number
LimitNOFILE=65535

# Specifies the maximum file size
LimitFSIZE=infinity

# Disable timeout
TimeoutStopSec=0

# SIGTERM signal is used to stop
KillSignal=SIGTERM
SendSIGKILL=no
SuccessExitStatus=143

# No timeout needed for simple type
# TimeoutStartSec=75

[Install]
WantedBy=multi-user.target
